# systemd/howtoindex.timer
[Unit]
Description=Run recollindex every hour

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
